# script/etc/crontabs/root
*/1 * * * * /root/thermo.sh >> /var/log/crond-root.log
